# tools/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

#sleep 1
#/root/Desktop/frp/1.sh &
#sleep 50 
#/root/Desktop/frp/frpc -c /root/Desktop/frp/frpc.ini &

sleep 1
/home/minieye/env/dist/assets/frp.sh &

exit 0
